# deploy/tmond-wifi.service
[Unit]
Description=tmon temperature monitoring daemon (WiFi push)
After=network.target

[Service]
Type=simple
User=tmon
Group=tmon
WorkingDirectory=/etc/tmon
ExecStart=/var/lib/tmon/venv/bin/tmon tmon.toml --transport wifi
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
